# debian/usblockout.postrm
#!/bin/sh

set -e

if [ "$1" = "remove" ]; then
	deb-systemd-helper reload dbus.service || true
fi

#DEBHELPER#
