# debian/usblockout.postinst
#!/bin/sh

set -e

case $1 in
	configure)
		deb-systemd-helper reload dbus.service || true
	;;
	abort-upgrade|abort-remove|abort-deconfigure)
	;;
	*)
		echo "postinst called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

#DEBHELPER#
